# examples.under-development/draft-example.nsenter/whoami.container
[Container]
ContainerName=whoami.example.com
Image=docker.io/traefik/whoami
Network=mynet.network
NoNewPrivileges=true
ReadOnly=true
