# examples/example3/caddy.container
[Container]
Exec=/usr/bin/caddy run --config /etc/caddy/Caddyfile
Image=localhost/caddy
Volume=%h/Caddyfile:/etc/caddy/Caddyfile:Z
