# examples/example4/whoami.container
[Container]
ContainerName=whoami
Image=docker.io/traefik/whoami
Network=mynet.network
