# examples.under-development/example99/example99.container.in
[Unit]
Requires=user@${envsubst_uid}.service
After=user@${envsubst_uid}.service
RequiresMountsFor=/run/user/${envsubst_uid}/containers

AssertPathIsDirectory=/home/${envsubst_user}/caddy_etc
AssertPathExists=/home/${envsubst_user}/caddy_etc/Caddyfile

[Service]
User=${envsubst_user}
ExecReload=podman exec example99 /usr/bin/caddy reload --config /etc/caddy/Caddyfile --address unix//run/admin.sock --force

[Container]
ContainerName=example99
Exec=/usr/bin/caddy run --config /etc/caddy/Caddyfile
Image=docker.io/library/caddy
Network=systemd-example99-net
Notify=true
Volume=/home/${envsubst_user}/caddy_etc:/etc/caddy:Z
