# examples.under-development/example99/whoami.container
[Container]
ContainerName=whoami
Image=docker.io/traefik/whoami
Network=example99-net.network
