# examples/example2/nginx.container
[Container]
ContainerName=nginx
Image=docker.io/library/nginx
Network=mynet.network
