# examples/example2/nginx.container
[Container]
ContainerName=nginx.example.com
Image=docker.io/library/nginx
Network=mynet.network
