# examples/example2/nginx.container
[Container]
Image=docker.io/library/nginx
Network=mynet.network
ContainerName=nginx.example.com
